# rootfs/etc/services.d/cron/run
#!/usr/bin/execlineb -P

with-contenv
/cron.sh
