# rootfs/etc/services.d/php-fpm/run
#!/usr/bin/execlineb -P

with-contenv
/entrypoint.sh php-fpm
